# okay-robot.service
[Unit]
Description=PiCar-X "Okay Robot" Voice-Activated Control System
Documentation=https://docs.sunfounder.com/projects/picar-x-v20/en/latest/
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/pi
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /home/pi/okay_robot.py
ExecStop=/bin/kill -SIGTERM $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal+console
StandardError=journal+console
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/root
Environment=DISPLAY=:0

# Allow access to hardware (I2C, GPIO, audio)
SupplementaryGroups=i2c gpio audio video input

# Watchdog: restart if hung
WatchdogSec=120

# Resource limits
TimeoutStartSec=60
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
